# crontab
* * * * * root cd /api && /bin/bash /api/lolstatus.sh >> /var/log/cron.log
