# crontab
* * * * * root /bin/bash /api/lolstatus.sh >> /var/log/cron.log 2>&1
