# crontab
* * * * * root /bin/bash /api/lolstatus.sh
